# overlay/common/udev-v7/71-persistent-net.rules
SUBSYSTEM=="net", ACTION!="remove", DRIVERS!="brcmfmac", ATTR{type}!="803", ENV{NM_UNMANAGED}="1", RUN+="/usr/local/sbin/createmon.sh $name %n"
